# websupport-vps-deployment/systemd/jobbie-infra-metrics.timer
[Unit]
Description=JOBBIE VPS infra metrics sampler every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
Unit=jobbie-infra-metrics.service

[Install]
WantedBy=timers.target
